# packages/lakka/L4T/assets/93-dock_hotplug.rules
ACTION=="change" KERNEL=="extcon0" RUN+="/usr/bin/dock-hotplug"
ACTION=="add" KERNEL=="extcon0" RUN+="/usr/bin/dock-hotplug"
