# projects/L4T/devices/Switch/filesystem/etc/udev/rules.d/99-nv-l4t-usb-device-mode.rules
#
# Copyright (c) 2019, NVIDIA CORPORATION.  All rights reserved.
#

ACTION=="change" KERNEL=="android0" SUBSYSTEM=="android_usb" RUN+="/opt/nvidia/l4t-usb-device-mode/nv-l4t-usb-device-mode-state-change.sh"
