# projects/L4T/devices/Switch/packages/usb-gadget-scripts/assets/usb-gadget.service
[Unit]
Description=Configure Nintendo Switch Gadget Framework

[Service]
ExecStart=/bin/usb-gadget.sh

[Install]
WantedBy=multi-user.target
