# projects/L4T/devices/Switch/packages/usb-gadget-scripts/assets/usb-tty.service
[Unit]
Description=Start TTY on ttyGS0
After=usb-gadget.service

[Service]
Type=oneshot
WorkingDirectory=/storage
Environment="ENV=/etc/profile"
ExecStart=/bin/sh < /dev/ttyGS0 > /dev/ttyGS0
Restart=always
KillMode=process
IgnoreSIGPIPE=no
# bash ignores SIGTERM
KillSignal=SIGHUP

[Install]
WantedBy=multi-user.target

